# cron/eramba-yearly.timer
[Unit]
Description=Run eramaba-yearly.service every year

[Timer]
OnCalendar=yearly

[Install]
WantedBy=multi-user.target
