# cron/eramba-daily.timer
[Unit]
Description=Run eramaba-daily.service every day

[Timer]
OnCalendar=*-*-* 00:30:00

[Install]
WantedBy=multi-user.target
