# cron/eramba-hourly.timer
[Unit]
Description=Run eramaba-hourly.service every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
